# provisioning/forebrain/templates/netherwick-forebrain-discovery.service.j2
[Unit]
Description=Netherwick forebrain capability discovery
After=network.target

[Service]
Type=simple
User={{ netherwick_service_user }}
Group={{ netherwick_service_group }}
ExecStart=/usr/local/bin/pete-higher-brain advertise --config /etc/netherwick/forebrain.toml
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/etc/netherwick

[Install]
WantedBy=multi-user.target
